# 06_systemd/30_seconds_log.timer
[Unit]
Description=Timer for logging every 30 seconds

[Timer]
OnUnitActiveSec=30s
OnBootSec=30s

[Install]
WantedBy=timers.target
