# 06_systemd/spawn-fcgi.service
[Unit]
Description=unit-version of spawn-fcgi config
After=local_fs.target

[Service]
Type=notify
EnvironmentFile=/etc/sysconfig/spawn-fcgi
ExecStart=/bin/spawn-fcgi start $OPTIONS
ExecStop=/bin/spawn-fcgi stop $OPTIONS
ExecReload=/bin/spawn-fcgi reload $OPTIONS
